# deploy/systemd/onepiece-rpg-backup.timer
[Unit]
Description=Backup diário One Piece RPG
[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true
RandomizedDelaySec=15m
[Install]
WantedBy=timers.target
