# deploy/systemd/onepiece-rpg-backup.service
[Unit]
Description=Backup One Piece RPG
[Service]
Type=oneshot
WorkingDirectory=/CAMINHO/DO/PROJETO
ExecStart=/usr/bin/make backup
